# scripts/kiosk.service
[Unit]
Description=GoldenMunch Kiosk Application
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
Environment="NODE_ENV=production"
WorkingDirectory=/home/pi/GoldenMunch_POS-System-With-Custom-Cake-Editor/client/Kiosk
ExecStartPre=/bin/sleep 10
ExecStart=/home/pi/GoldenMunch_POS-System-With-Custom-Cake-Editor/scripts/start-kiosk.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
